# recipes-rtd-generic/rtd-dpkg/files/info/rtd-services.list
/.
/etc
/etc/node_red
/etc/node_red/nodered
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/ble.service
/lib/systemd/system/restart_ble.service
/lib/systemd/system/sshd.service
/lib/systemd/system/watchdog.service
/lib/systemd/system/wifi_enable.service
/lib/systemd/system-preset
/lib/systemd/system-preset/98-rtd-services.preset
/usr
/usr/bin
/usr/bin/ble.sh
/usr/bin/restart_ble_serial.sh
/usr/bin/watchdog
/usr/bin/wifi_enable.sh
