# recipes-rtd-generic/rtd-dpkg/files/info/rtd-mcu-bin.prerm
#!/bin/sh
[ "$1" != "upgrade" ] || exit 0
set -e
if type systemctl >/dev/null 2>/dev/null; then
	if [ -z "$D" ]; then
		systemctl stop mcu_sbl_mode.service

		systemctl disable mcu_sbl_mode.service
	fi
fi
